# crontab
* * 16 ? * * * node supercronic/tengah_malam.js > log.txt
* * 19 ? * * * node supercronic/tiga_pagi.js > log.txt